# config/systemd/user/onedrive-sync.timer
[Unit]
Description=Run OneDrive synchronization hourly

[Timer]
OnBootSec=15min
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
