# config/systemd/user/onedrive-sync.service
[Unit]
Description=Syncs OneDrive

[Service]
ExecStart=/usr/bin/onedrive --synchronize

[Install]
WantedBy=default.target
